# contrib/vomacs.service
[Unit]
Description=vomacs voice transcription daemon
After=graphical-session.target
PartOf=graphical-session.target
Wants=vomacs-kde-helper.service

[Service]
Type=simple
ExecStart=/usr/bin/vomacs serve
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
